# dot_config/i3/crontab
*/5 * * * * curl -s wttr.in/Andradina?format=4 | head -n 7 > ~/.weather.cache
